# etc/supervisor/conf.d/historyserver.conf
[program:historyserver]
user=hadoop
command=/opt/hadoop-2.7.0/bin/mapred historyserver